# templates/etc/init/gor.conf.erb
description "Gor traffic replay"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 20

<%  @envvars.each do |k, v| %>
<%=   "env #{k}='#{v}'" -%>
<%  end %>

exec <%= @binary_path %> \
  <%=
@args.sort.map { |k,v|
  [v].flatten.map { |i|
    if i == ''
      "#{k}"
    else
      "#{k}='#{i}'"
    end
  }
}.join(" \\\n  ")
%>
