# templates/etc/init/gor.conf.erb
description "Gor traffic replay"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 20

exec /usr/bin/gor \
  <%=
@args.sort.map { |k,v|
  [v].flatten.map { |i|
    "#{k}='#{i}'"
  }
}.join(" \\\n  ")
%>
