# templates/etc/init/gor.conf.erb
description "Gor traffic replay"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 20

exec /usr/bin/gor \
  <%= @args.map { |k,v| "#{k}='#{v}'" }.join(" \\\n  ") %>
